# packaging/launchd/com.clawnify.agent-control-ui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Launchd agent for @clawnify/agent-control-ui — auto-starts the local
  static server on a Mac mini or any macOS host running OpenClaw.

  Install:
    cp packaging/launchd/com.clawnify.agent-control-ui.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.clawnify.agent-control-ui.plist

  Stop / remove:
    launchctl unload ~/Library/LaunchAgents/com.clawnify.agent-control-ui.plist

  Logs:
    tail -f /tmp/agent-control-ui.{out,err}.log

  Adjust the ProgramArguments path if `npx` lives somewhere other than
  /opt/homebrew/bin (Apple silicon) — /usr/local/bin on Intel Macs.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clawnify.agent-control-ui</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npx</string>
    <string>--yes</string>
    <string>@clawnify/agent-control-ui</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>5174</string>
    <key>HOST</key>
    <string>127.0.0.1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/agent-control-ui.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/agent-control-ui.err.log</string>
</dict>
</plist>
